# deploy/systemd/github-ai-reviewer.service
[Unit]
Description=GitHub AI Reviewer
Documentation=https://github.com/icatw/github-ai-reviewer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/github-ai-reviewer
EnvironmentFile=/home/ubuntu/github-ai-reviewer/.env.production
ExecStart=/home/ubuntu/github-ai-reviewer/server
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/ubuntu/github-ai-reviewer

[Install]
WantedBy=multi-user.target
